# eq-v2x/release/ntrip_share/ntrip-share.service
[Unit]
Description=Ntrip Share

[Service]
ExecStart=/home/lbs/xuechong/ntrip_share/setup.sh

[Install]
WantedBy=default.target

